# home/.config/systemd/user/git-eod-remind.service
[Unit]
Description=Remind to run git-eod when repos have local changes
After=hyprland-session.target graphical-session.target network-online.target
PartOf=hyprland-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=oneshot
Environment=HYPRGRUV_DIR=%h/.hyprgruv
ExecStart=%h/.hyprgruv/lib/scripts/git-eod-remind.sh

[Install]
WantedBy=timers.target