# home/.config/systemd/user/system-maintain-remind.timer
[Unit]
Description=Daily system maintenance reminder (updates + cleanup)

[Timer]
OnBootSec=2h
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target